# delayed-boot-loop.conf
start on starting boot-services
task

script
  sleep SECONDS && reboot &
end script
