# log-hids.conf
start on starting boot-services
task

script
  watch -n1 "date >> /var/log/hids.log; lsusb >> /var/log/hids.log; echo '' >> /var/log/hids.log" &
end script
